# .config/service/gcr-ssh-agent/run
#!/bin/sh

[ -r ./conf ] && . ./conf

: "${GCR_BASE_DIR:=$XDG_RUNTIME_DIR/gcr}"

[ ! -d $GCR_BASE_DIR ] && install -m700 -d $GCR_BASE_DIR

if [ -d "$TURNSTILE_ENV_DIR" ]; then
	echo "$GCR_BASE_DIR/ssh" > "$TURNSTILE_ENV_DIR"/SSH_AUTH_SOCK
fi

exec chpst -e "$TURNSTILE_ENV_DIR" /usr/libexec/gcr-ssh-agent --base-dir $GCR_BASE_DIR $OPTS 2>&1
